# docker-ubuntu-desktop/supervisord.conf
[program:vnc]
priority=20
directory=/
command=/opt/vnc_start.sh
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/tightvnc.log
redirect_stderr=true